# curator4-Irina.conf
[program:curator4-Irina]
command=python3 /home/curator4/curator4-Irina/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-Irina/test.err.log
stdout_logfile=/home/curator4/curator4-Irina/test.out.log